# debian/orchestrator/DEBIAN/postinst
#!/bin/sh

systemctl enable arduino-orchestrator

if [ -d /run/systemd/system ]; then
  systemctl daemon-reload
  systemctl start arduino-orchestrator
fi
