# debian/orchestrator/DEBIAN/prerm
#!/bin/sh

systemctl disable arduino-orchestrator

if [ -d /run/systemd/system ]; then
  systemctl stop arduino-orchestrator
fi
