# debian/arduino-app-cli/DEBIAN/postinst
#!/bin/sh

chown -R arduino:arduino /home/arduino/.config/arduino-app-cli/examples

systemctl enable arduino-app-cli
