# debian/arduino-app-cli/DEBIAN/prerm
#!/bin/sh

systemctl disable arduino-app-cli
